# deploy/ngrok.service
[Unit]
Description=ngrok Tunnel — https://geserverhub.ngrok.app → port 3005
After=network.target geserver-frontend.service
Wants=geserver-frontend.service

[Service]
Type=simple
User=pavinee
ExecStart=/usr/local/bin/ngrok start --all --config /home/pavinee/.config/ngrok/ngrok.yml
Restart=alway
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
